# deploy/intraday-bot.service
# Alpaca Paper Trading Bot — Intraday systemd service unit
#
# Runs the intraday momentum/mean-reversion scanner and execution engine.
# Requires AUTO_EXECUTE=true in .env for automatic order placement.
#
# Installation (run as root on the VPS):
#   cp intraday-bot.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable intraday-bot
#   systemctl start  intraday-bot
#
# Useful commands:
#   systemctl status  intraday-bot
#   systemctl stop    intraday-bot
#   journalctl -u intraday-bot -f

[Unit]
Description=Intraday Trading Bot (Alpaca)
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/trading-bot

EnvironmentFile=/opt/trading-bot/.env

ExecStart=/opt/trading-bot/venv/bin/python main.py

# ── Restart policy ──────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=30s

# ── Logging ─────────────────────────────────────────────────────────────────
# Let Python logging_config.py own file logging under /opt/trading-bot/logs.
# Keeping stdout/stderr in journald avoids root-owned file creation conflicts.
StandardOutput=journal
StandardError=journal

# ── Hardening (optional but recommended) ────────────────────────────────────
# PrivateTmp=yes
# NoNewPrivileges=yes
# ProtectSystem=strict
# ReadWritePaths=/opt/trading-bot/logs /opt/trading-bot/data

[Install]
WantedBy=multi-user.target
